# jetson/systemd/user-vision-landing.service
[Unit]
Description=UAV Vision Landing Runtime (User)
After=network-online.target

[Service]
Type=simple
ExecStart=/bin/bash -lc 'cd "$HOME/on-board computer/jetson" && export PYTHONPATH="$PWD/src:${PYTHONPATH:-}" && exec /usr/bin/python3 -m vision_landing.main --config configs/aruco_live.yaml --camera-config configs/camera.yaml'
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
